# android-mtp@.service
# Android MTP mount service

[Service]
Type=forking
ExecStart=/usr/bin/mtp mount %I
ExecStop=/usr/bin/mtp unmount %I
